# pi-media-station.service
# Systemd Service für Auto-Start beim Boot
# Datei: /etc/systemd/system/pi-media-station.service

[Unit]
Description=Raspberry Pi Media Station
After=graphical-session.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/pi_media_station
Environment=DISPLAY=:0
ExecStart=/usr/bin/python3 /home/pi/pi_media_station/main.py
Restart=always
RestartSec=10

[Install]
WantedBy=graphical-session.target

# Installation:
# sudo cp pi-media-station.service /etc/systemd/system/
# sudo systemctl daemon-reload
# sudo systemctl enable pi-media-station.service
# sudo systemctl start pi-media-station.service
